# deploy/oracle/levtrade-collector.service
[Unit]
Description=LevTrade collector loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/levtrade/collector/app
EnvironmentFile=/srv/levtrade/collector/.env.collector
ExecStart=/usr/bin/npm run collector:start
Restart=always
RestartSec=10
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
